# systemd/triggerhappy.socket
[Socket]
ListenDatagram=%t/thd.socket

[Install]
WantedBy=sockets.target
